# etc/services.d/litestream/run
#!/usr/bin/with-contenv bash
/usr/local/bin/litestream replicate
