# etc/services.d/litestream/run
#!/usr/bin/with-contenv bash
/usr/local/bin/litestream replicate /var/lib/gha/db s3://dev.gha.fly.dev/db
