# etc/services.d/gha/run
#!/usr/bin/with-contenv bash
/usr/local/bin/gha -ingest-rate 10 /data/db 
