# etc/services.d/gha/run
#!/usr/bin/execlineb -P
/usr/local/bin/gha -ingest-rate 10 /var/lib/gha/db
